# deploy/systemd/unfollow-tracker-fetcher.service
[Unit]
Description=Unfollow Tracker Fetcher Worker
After=network.target

[Service]
Type=simple
User=unfollowtracker
WorkingDirectory=/opt/unfollow-tracker
ExecStart=/opt/unfollow-tracker/fetcher
EnvironmentFile=/opt/unfollow-tracker/.env
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
